# fdns/fwd_dns/src/fwd_dns_bin/crontab/.cronjob
0 * * * * /bin/sh @base/crontab/get_topn_data.sh
