# adns/adns_tsar/conf/adns_tsar.cron
# cron adns_tsar collect once per minute
MAILTO=""
* * * * * root /usr/bin/adns_tsar --cron > /dev/null 2>&1
